# raspotify/debian/postrm
#!/bin/sh

# Based off what debhelper would have added

set -e

if [ "$1" = "remove" ]; then
    if [ -x "/usr/bin/deb-systemd-helper" ]; then
        deb-systemd-helper mask raspotify.service >/dev/null
    fi
fi

if [ "$1" = "purge" ]; then
    if [ -x "/usr/bin/deb-systemd-helper" ]; then
        deb-systemd-helper purge raspotify.service >/dev/null
        deb-systemd-helper unmask raspotify.service >/dev/null
    fi
fi

if [ -d /run/systemd/system ]; then
    systemctl --system daemon-reload >/dev/null || true
fi
