# raspotify/lib/systemd/system/raspotify-crash-report-generator.service
[Unit]
Description=Raspotify Crash Report Generator

[Service]
Type=oneshot
ExecStart=/usr/bin/raspotify-crash-report-generator.sh
